# my-config/etc/systemd/user/mpDris.service.d/unit.conf
[Unit]
ConditionPathExists=/usr/bin/mpd
ConditionPathExists=/usr/bin/mpDris2

