# my-config/etc/systemd/user/mpd.service.d/environ.conf
[Service]
Environment=MPD_HOST=%t/mpd.socket
Environment=MPD_PORT=6%U
