# my-config/etc/systemd/user/mpd.service.d/service.conf
[Unit]
ConditionPathExists=%h/.config/mpd/mpd.conf
ConditionPathExists=/usr/bin/mpd
